# apps/rocana/build/cookbooks/rocana/templates/rocana-metadata-consumer.service.erb
[Unit]
Description=Rocana Metadata Consumer
After=network.target


[Service]
EnvironmentFile=/opt/gex/rocana-env
ExecStart=/usr/bin/nohup <%=@rocana_base_dir%>bin/rocana-metadata-consumer
Restart=always
User=rocana
TimeoutStartSec=300
TimeoutStopSec=300



[Install]
WantedBy=multi-user.target

